# meta-smart-house/recipes-connectivity/gpio-client/files/gpio.service
[Unit]
Description=Startup script for gpio client/doors

[Service]
ExecStart=/usr/bin/gpio-doors

[Install]
WantedBy=multi-user.target